# playbooks/nexus.service.j2
[Unit]
Description=Nexus Repository Manager
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/nexus start
ExecStop=/usr/bin/nexus stop
User=nexus
Group=nexus
LimitNOFILE=65536
Restart=on-failure

[Install]
WantedBy=multi-user.target